# freePBX/asterisk/supervisor/conf.d/mysql.conf
[program:mysql]
command=service mysqld start
#/etc/service/mysqld/run
# autorestart=true
# autorestart=false
# log_stderr=true 
# startretries=0
priority=3332
# exitcodes=0
# user=root


